# env/config/supervisord/supervisord.conf
[supervisord]
logfile=/var/log/supervisor.log ;
logfile_maxbytes=50MB        ;
logfile_backups=10           ;
loglevel=info                ;
pidfile=/tmp/supervisord.pid ;
nodaemon=false               ;
minfds=1024                  ;
minprocs=200                 ;


[program:hangman]
command=/usr/local/bin/gunicorn conf.wsgi:application --bind 127.0.0.1:7007 --pid /tmp/app.pid ;
directory=/deploy/apps/hangman-game/app/ ;


[inet_http_server]
port=127.0.0.1:9001   ;

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=http://127.0.0.1:9001 ;
